# roles/base/files/etc/cron.d/custom
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# m  h  dom mon dow user command
*/2  *  *   *   *   root /usr/local/bin/ntfy-health-check >> /var/log/ntfy-health-check.log
58   23 *   *   *   root /usr/local/bin/ntfy-stats > /dev/null
*/30 *  *   *   *   root /usr/local/bin/ntfy-cachedb-backup > /dev/null 2>&1
*/4  *  *   *   *   root /usr/local/bin/checkdiskspace
