# roles/base/files/etc/cron.d/custom
# m h  dom mon dow   command

*/2   * * * * /usr/local/bin/ntfy-health-check >> /var/log/ntfy-health-check.log
58   23 * * * /usr/local/bin/ntfy-stats > /dev/null
*/30  * * * * /usr/local/bin/ntfy-cachedb-backup > /dev/null 2>&1
*/4   * * * * /usr/local/bin/checkdiskspace
